# mkdocs/rootfs/etc/services.d/mkdocs/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Filii Lamberti Hass.io Add-ons: MkDocs
# Runs MkDocs
# ==============================================================================

# Go to mkdocs directory
cd /share/mkdocs

# Build the docs
mkdocs build --clean

# Go to site directory
cd site

# And finally serve it
exec python3 -m http.server 8099
